# deploy/sap-facture.service
[Unit]
Description=SAP-Facture Platform - Invoice Management System
Documentation=https://github.com/your-org/sap-facture
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=sap
Group=sap
WorkingDirectory=/opt/sap-facture

# Pre-start checks
ExecStartPre=/usr/bin/docker --version

# Start the service
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up

# Stop the service
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down

# Restart policy
Restart=always
RestartSec=10s
StartLimitInterval=60s
StartLimitBurst=3

# Resource limits
MemoryLimit=2G
CPUQuota=200%

# Security
PrivateTmp=yes
NoNewPrivileges=yes
ReadWritePaths=/opt/sap-facture/data /opt/sap-facture/storage /opt/sap-facture/logs

# Timeouts
TimeoutStartSec=60s
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
